# systemd/glorytun@.service.in
[Unit]
Description=Glorytun on %I
After=network.target

[Service]
Type=simple
Restart=always
EnvironmentFile=/etc/glorytun/%i/env
ExecStart=@bindir@/glorytun-run keyfile /etc/glorytun/%i/key $OPTIONS
ExecStartPost=-/etc/glorytun/%i/post.sh
CapabilityBoundingSet=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
